# lamp/php7.4/lamp.supervisord.conf
[program:php-fpm]
command=php-fpm7.4 -F -c /etc/php/7.4/fpm/php-fpm.conf

[program:apache2]
command=/usr/sbin/apachectl -DFOREGROUND
