# ubuntu/supervisord.conf
[supervisord]
nodaemon=true
user=root

[supervisorctl]

[program:startup]
command=/startup.sh
