# lamp/template/lamp.supervisord.conf
[program:php-fpm]
command=php-fpm{{php_version}} -F -c /etc/php/{{php_version}}/fpm/php-fpm.conf

[program:apache2]
command=/usr/sbin/apachectl -DFOREGROUND
