# lamp/php8.1/lamp.supervisord.conf
[program:php-fpm]
command=php-fpm8.1 -F -c /etc/php/8.1/fpm/php-fpm.conf

[program:apache2]
command=/usr/sbin/apachectl -DFOREGROUND
