# ubuntu/template/supervisord.conf
[supervisord]
nodaemon=true
user=root

[supervisorctl]

[program:sshd]
command=/usr/sbin/sshd -D

[program:startup]
command=/startup.sh
startsecs=0
autorestart=false
