# ubuntu/template/supervisord.conf
[supervisord]
nodaemon=true
user=root

[supervisorctl]

[program:startup]
command=/startup.sh
startsecs=0
autorestart=false
